# systemd/trainer-scheduler.service
[Unit]
Description=GH200 Model Trainer Scheduler
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/gh200-trading-system
ExecStart=/usr/bin/python3 /home/ubuntu/gh200-trading-system/monitoring/schedule/trainer_scheduler.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
