# scripts/trading-scheduler.service
[Unit]
Description=Trading System Scheduler
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/inavvi2
ExecStart=/usr/bin/python3 /home/ubuntu/inavvi2/scripts/trading_scheduler.py
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=trading-scheduler
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target